# priv/install/templates/etc/supervisord/staging.conf
[program:<%= application_name %>_staging]
environment=MIX_ENV=staging
command=mix phoenix.server
directory=/sites/staging/<%= application_name %>
user=<%= application_name %>
autostart=false
autorestart=true
stdout_logfile=/sites/staging/<%= application_name %>/logs/supervisord.log
redirect_stderr=true